# systemd-scripts/backup-example.timer
[Unit]
Description=Timer for example backup script

[Timer]
OnCalendar=02:00
Unit=backup-example.service

[Install]
WantedBy=multi-user.target
